# systemd/fotios-web.service
[Unit]
Description=Fotios Claude Web Panel (Flask)
After=network.target mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/codehero/web
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/codehero/web/app.py
ExecStopPost=/bin/bash -c 'fuser -k 5000/tcp 2>/dev/null || true'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
